# files/etc/cron.d/statsclean
## PUPPETIZED

10 1 * * *	root	/root/scripts/cleanStats.sh
